# resources/profile.d/profile_d_composer.sh
export PATH=$PATH:COMPOSER_GLOBAL_DIR/vendor/bin/
